# scripts/launchd/com.jarvis.devops.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.devops</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS/scripts/devops_gcp_monitor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS/logs/devops_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS/logs/devops_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
