# scripts/launchd/com.jarvis.keeper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.keeper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS/scripts/keeper_reminders.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS/logs/keeper_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jitheesh.pj/Desktop/Jitheesh/JARVIS/logs/keeper_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
